# resource/etc/cron.crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/bin:/usr/sbin
HOME=/
* * * * * /bin/bash /root/manage_flash_cap.sh > /opt/private/log/spi-flash.recorder
